# ansible/roles/imapnotify/files/imapnotify.service
[Unit]
Description=imapnotify service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/imapnotify -c %h/.config/imapnotify/config.js
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
